# hms-assist.service
[Unit]
Description=HMS-Assist Voice Assistant API v2
After=network.target postgresql.service hms-assist-sync.service
Wants=network.target

[Service]
Type=simple
User=aamat
WorkingDirectory=/var/lib/hms-assist
ExecStart=/usr/local/bin/hms_assist
Environment="HMS_ASSIST_CONFIG=/etc/hms-assist/config.yaml"

Restart=always
RestartSec=10
NoNewPrivileges=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hms-assist

[Install]
WantedBy=multi-user.target
